# environment_bak/app/profile.d/env.sh
export DB_HOST=mongodb://192.168.10.150/posts